# systemd/99-amdgpu-runtime-pm-on.rules
ACTION=="add|bind", SUBSYSTEM=="pci", DRIVER=="amdgpu", ATTR{power/control}="on", RUN+="/usr/local/sbin/amdgpu-runtime-pm-guard --apply"
