# systemd/localllm-gpu-power-limit.service
[Unit]
Description=Apply NVIDIA GPU power limit for localllm host
After=systemd-modules-load.service
ConditionPathExists=/usr/bin/nvidia-smi

[Service]
Type=oneshot
Environment=GPU_POWER_LIMIT_WATTS=300
EnvironmentFile=-/etc/default/localllm-gpu-power-limit
ExecStart=/usr/bin/nvidia-smi -pl ${GPU_POWER_LIMIT_WATTS}

[Install]
WantedBy=multi-user.target
