# systemd/localllm-main.service
[Unit]
Description=localllm main llama.cpp service

[Service]
Type=simple
WorkingDirectory=%h/projects/localllm
Environment=MAIN_HOST=127.0.0.1
Environment=MAIN_PORT=8091
EnvironmentFile=%h/projects/localllm/config/localllm-main.env
ExecStart=%h/projects/localllm/scripts/serve-main.sh
Restart=always
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=default.target
